# .config/systemd/user/notmuch.timer
[Unit]
Description=notmuch timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=default.target
