# .config/systemd/user/notmuch.service
[Unit]
Description=Update notmuch

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new --quiet
Restart=no

[Install]
WantedBy=default.target
