# .config/systemd/user/xsession.target
[Unit]
Description=Xsession managed by systemd
BindsTo=graphical-session.target
