# .config/systemd/user/systemd-tmpfiles-clean.timer.d/override.conf
[Timer]
OnUnitActiveSec=
OnUnitActiveSec=1h
